# example.cron
#!/bin/sh
# put in /etc/cron.daily
# this will rebuild certificate (as root) if it got updated by acmetool (crontab -u acme)
make --silent -C /opt/zcs-acmetool domain.crt
